# egress/reverse/systemd/dxreverse-hub.service
[Unit]
Description=Daxiang Android reverse egress hub
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/daxiang/dxreverse/dxreverse server --config /etc/daxiang/dxreverse/server.yaml
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/log/daxiang /etc/daxiang/dxreverse

[Install]
WantedBy=multi-user.target
